# scripts/systemd/newfuhi.service
[Unit]
Description=NewFUHI Gunicorn
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/NewFUHI
EnvironmentFile=/home/ubuntu/NewFUHI/.env
Environment=ENV=production
Environment=DJANGO_ENVIRONMENT=production
Environment=DJANGO_SETTINGS_MODULE=project.settings
ExecStart=/home/ubuntu/NewFUHI/.venv/bin/gunicorn project.wsgi:application \
    --bind 127.0.0.1:8000 --workers 2 --timeout 30 \
    --access-logfile /var/log/newfuhi/gunicorn-access.log \
    --error-logfile /var/log/newfuhi/gunicorn-error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
